# deploy/studentstay.service
[Unit]
Description=StudentStay Node.js application
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/telebe-hotel-main/server
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=3002
Environment=ADMIN_USER=admin
Environment=ADMIN_PASSWORD_HASH=scrypt:studentstayadmin2026:e2317b3a2defe0838e88b596d759eb64788da1dce757d46b4f282dbadada0a8812a7120b68a9f3560b91288c906d32b8c1145c54a486c9bddf7ef3cdc5498209
Environment=ADMIN_NOTIFY_EMAIL=

[Install]
WantedBy=multi-user.target
